# systemd/scanbox-stack.service
[Unit]
Description=SCANBOX Docker Compose Stack
# Must start after Docker is up and the USB gadget is configured.
Requires=docker.service
After=docker.service scanbox-gadget.service
# Ensure network is available (usb0 needs it for DHCP container).
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/Alfred/scanbox

# 'docker compose up -d' honours depends_on ordering:
#   scanbox_dhcp → immediately
#   vid_mux_test → immediately (compiles kernel module, creates /dev/video200)
#   vid_mux      → waits until vid_mux_test healthcheck passes
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose stop

# Give vid_mux_test time to compile v4l2loopback (can take 60-90s on first boot).
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
